# massnmap.timer
[Unit]
Description=RSS Fetcher Timer

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
