# massnmap.timer
[Unit]
Description=RSS Fetcher Timer

[Timer]
OnCalendar=3w 
Persistent=true

[Install]
WantedBy=timers.target
